# scripts/com.contentmasterpro.ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contentmasterpro.ingest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jonathanedwards/AUTOMATION/SubStack/content-master-pro/scripts/daily-ingest.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jonathanedwards/AUTOMATION/SubStack/content-master-pro/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jonathanedwards/AUTOMATION/SubStack/content-master-pro/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
